# app/hadoop-0.20.2-hdfs/setaside-hdfs.socket
[Unit]
Description=Set-aside Hadoop Distributed File System Socket

[Socket]
ListenStream=54310
Accept=false

[Install]
