# conf/echoserver-new.socket
[Unit]
Description=EchoServer - New style socket unit

[Socket]
ListenStream=10008
Accept=false

[Install]
