# conf/echoselectorserver-new.socket
[Unit]
Description=EchoSelectorServer - New style socket unit

[Socket]
ListenStream=10009
Accept=false

[Install]
